# deploy/postcraft.service
[Unit]
Description=PostCraft API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=@POSTCRAFT_ROOT@/backend
EnvironmentFile=@POSTCRAFT_ROOT@/.env
ExecStart=@POSTCRAFT_ROOT@/backend/.venv/bin/python -m uvicorn app.main:app \
  --host 127.0.0.1 --port 18231 --workers 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
